# scripts/zhuli/launchd/com.howard.zhuli.full_scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.howard.zhuli.full_scanner</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/howard/.pyenv/shims:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
		<key>FINMIND_TOKEN</key>
		<string>eyJ0eXAiOiJKV1QiLCJhbGciOiJIUzI1NiJ9.eyJkYXRlIjoiMjAyNi0wMy0yMyAyMjowODoyOSIsInVzZXJfaWQiOiJwYzAyODc3MSIsImVtYWlsIjoicGMwMjg3NzFAZ21haWwuY29tIiwiaXAiOiIxMTguMTY4LjUyLjEyNiJ9.FDP59dSzQ5zCIAsV5w9dpQ48bWeqJsATYheMFgs_NGk</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>mkdir -p /tmp/launchd; echo "=== START $(date +%Y-%m-%d_%H:%M:%S) ==="; echo "--- step 1/2: daily_scanner_job ---"; /Users/howard/.pyenv/shims/python3 /Users/howard/Repository/stock-k-bar/scripts/zhuli/daily_scanner_job.py; echo "--- step 2/2: overnight_swing scanner ---"; /Users/howard/.pyenv/shims/python3 /Users/howard/Repository/stock-k-bar/scripts/zhuli_scanner.py --signal overnight_swing --date $(date +%Y-%m-%d); echo "=== END $(date +%Y-%m-%d_%H:%M:%S) exit=$? ==="</string>
	</array>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>2</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>3</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>5</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
	</array>
	<key>StandardOutPath</key>
	<string>/tmp/launchd/com.howard.zhuli.full_scanner.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/launchd/com.howard.zhuli.full_scanner.err</string>
	<key>WorkingDirectory</key>
	<string>/Users/howard/Repository/stock-k-bar</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
