# scripts/zhuli/launchd/com.howard.zhuli.closed_scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.howard.zhuli.closed_scanner</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FINMIND_TOKEN</key>
        <string>eyJ0eXAiOiJKV1QiLCJhbGciOiJIUzI1NiJ9.eyJ1c2VyX2lkIjoicGMwMjg3NzEiLCJlbWFpbCI6InBjMDI4NzcxQGhvdG1haWwuY29tIiwidG9rZW5fdmVyc2lvbiI6MH0.v-g0kPNYRs3pWQMk_UzjvGYZUcitFnEhRDSZAbpDnHo</string>
        <key>FUBON_API_KEY</key>
        <string>CAD10AD5F4F0EA880FC3529FC7C7857D08228988EE9F718B7E83A8EE8C3B196A</string>
        <key>FUBON_CREDENTIAL_FILE</key>
        <string>/Users/howard/.fubon/A123309212_20270401.p12</string>
        <key>FUBON_CREDENTIAL_PWD</key>
        <string>My253845</string>
        <key>FUBON_PID</key>
        <string>A123309212</string>
        <key>FUBON_PWD</key>
        <string>camGep2zaz</string>
        <key>PATH</key>
        <string>/Users/howard/.pyenv/versions/3.12.12/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>mkdir -p /tmp/launchd
echo "=== START $(date +%Y-%m-%d_%H:%M:%S) ==="
echo "--- step 1/2: sync_today (Fubon snapshot → DB) ---"
/Users/howard/.pyenv/versions/3.12.12/bin/python3 /Users/howard/Repository/stock-k-bar/scripts/zhuli/sync_today.py --skip-institutional
echo "--- step 2/2: overnight_swing scanner ---"
/Users/howard/.pyenv/versions/3.12.12/bin/python3 /Users/howard/Repository/stock-k-bar/scripts/zhuli_scanner.py --signal overnight_swing --date $(date +%Y-%m-%d)
echo "=== END $(date +%Y-%m-%d_%H:%M:%S) exit=$? ==="</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/howard/Repository/stock-k-bar</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>35</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/launchd/com.howard.zhuli.closed_scanner.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/launchd/com.howard.zhuli.closed_scanner.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
