# scripts/zhuli/launchd/com.howard.zhuli.evening_brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.howard.zhuli.evening_brief</string>

    <!--
      每晚 18:00 深度報告（FinMind 分點資料 ~17:00 後 ready）：
        - 大哥資料（即時抓取、不用 cache）
        - holdings 狀態分析
        - scanner 結果摘要
        - trade_journal 當日
        - 明日 SOP 候選
      輸出：docs/主力大課程/daily_brief/YYYY-MM-DD_evening.md
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/howard/.pyenv/shims/python3</string>
        <string>/Users/howard/Repository/stock-k-bar/scripts/zhuli/evening_brief.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/howard/.pyenv/shims:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>FINMIND_TOKEN</key>
        <string>eyJ0eXAiOiJKV1QiLCJhbGciOiJIUzI1NiJ9.eyJ1c2VyX2lkIjoicGMwMjg3NzEiLCJlbWFpbCI6InBjMDI4NzcxQGhvdG1haWwuY29tIiwidG9rZW5fdmVyc2lvbiI6MH0.v-g0kPNYRs3pWQMk_UzjvGYZUcitFnEhRDSZAbpDnHo</string>
    </dict>

    <!-- 每週一~五 21:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/zhuli_evening_brief.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/zhuli_evening_brief.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
