# scripts/zhuli/launchd/com.howard.zhuli.evening_fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.howard.zhuli.evening_fetch</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/howard/.pyenv/shims:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
		<key>FINMIND_TOKEN</key>
		<string>eyJ0eXAiOiJKV1QiLCJhbGciOiJIUzI1NiJ9.eyJ1c2VyX2lkIjoicGMwMjg3NzEiLCJlbWFpbCI6InBjMDI4NzcxQGhvdG1haWwuY29tIiwidG9rZW5fdmVyc2lvbiI6MH0.v-g0kPNYRs3pWQMk_UzjvGYZUcitFnEhRDSZAbpDnHo</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>mkdir -p /tmp/launchd; echo "=== START $(date +%Y-%m-%d_%H:%M:%S) ==="; /bin/bash /Users/howard/Repository/stock-k-bar/scripts/zhuli/evening_fetch.sh; echo "=== END $(date +%Y-%m-%d_%H:%M:%S) exit=$? ==="</string>
	</array>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>2</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>3</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>5</integer>
			<key>Hour</key>
			<integer>21</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</array>
	<key>StandardOutPath</key>
	<string>/tmp/launchd/com.howard.zhuli.evening_fetch.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/launchd/com.howard.zhuli.evening_fetch.err</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
