# scripts/zhuli/launchd/com.howard.zhuli.daily_scanner_evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.howard.zhuli.daily_scanner_evening</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FINMIND_TOKEN</key>
        <string>eyJ0eXAiOiJKV1QiLCJhbGciOiJIUzI1NiJ9.eyJkYXRlIjoiMjAyNi0wMy0yMyAyMjowODoyOSIsInVzZXJfaWQiOiJwYzAyODc3MSIsImVtYWlsIjoicGMwMjg3NzFAZ21haWwuY29tIiwiaXAiOiIxMTguMTY4LjUyLjEyNiJ9.FDP59dSzQ5zCIAsV5w9dpQ48bWeqJsATYheMFgs_NGk</string>
        <key>PATH</key>
        <string>/Users/howard/.pyenv/shims:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>mkdir -p /tmp/launchd
echo "=== START $(date +%Y-%m-%d_%H:%M:%S) ==="
echo "--- step 1/2: daily_scanner_job (結構 scanner) ---"
/Users/howard/.pyenv/shims/python3 /Users/howard/Repository/stock-k-bar/scripts/zhuli/daily_scanner_job.py --enable-kline-tier-a
DSE_EXIT=$?
echo "--- step 2/2: 6 個策略 scanner (每個獨立輸出 CSV) ---"
for SIG in suffocation institutional_firstbuy swing_breakout reversal_breakout pennant_flag institutional_swing; do
    echo "  [scanner: $SIG]"
    /Users/howard/.pyenv/shims/python3 /Users/howard/Repository/stock-k-bar/scripts/zhuli_scanner.py --signal "$SIG" --date $(date +%Y-%m-%d) 2&gt;&amp;1 | tail -3
done
bash /Users/howard/Repository/stock-k-bar/scripts/zhuli/notify_kline_tier_a.sh
echo "=== END $(date +%Y-%m-%d_%H:%M:%S) exit=$DSE_EXIT ==="</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/howard/Repository/stock-k-bar</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>45</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>45</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>45</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>45</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>45</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/launchd/com.howard.zhuli.daily_scanner_evening.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/launchd/com.howard.zhuli.daily_scanner_evening.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
